# app/config/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:flask]
command=/usr/local/bin/flask run
directory=/app
